# docker/web/etc/cron.d/local.conf
*/3 * * * * cd /var/www/html/app && composer dumpautoload >> /dev/null 2>&1
