# docker/web/etc/cron.d/production.conf
0 * * * * docker exec -i prisma-service-web php /var/www/html/api/view/outros/rotinasgeral.php >> /var/log/rotinasgeral.log 2>&1
